# cookbooks/cpp_hydra/templates/default/cpp_hydra.init.d.erb
#!/sbin/runscript

start() {
  ebegin "Starting Hydra..."
  cd /data/cpp_hydra/current && /usr/bin/env DAEMON_ENV=<%= @environment %> bundle exec /data/cpp_hydra/current/bin/cpp_hydra start
  eend $?
}

stop() {
  ebegin "Stopping Hydra..."
  cd /data/cpp_hydra/current && /usr/bin/env DAEMON_ENV=<%= @environment %> bundle exec /data/cpp_hydra/current/bin/cpp_hydra stop 
  eend $?
}

restart() {
  ebegin "Restarting Hydra..."
  stop
  start
  eend $?
}
